# templates/ec2-environment.service.j2
[Unit]
Description=AWS EC2 Environment
After=systemd-networkd-wait-online.service
Wants=systemd-networkd-wait-online.service
After=systemd-networkd.service
Wants=systemd-networkd.service
Before=network-online.target

[Service]
Type=oneshot
ExecStart=retry /usr/bin/ec2-environment
Environment=RETRY_MAX_ATTEMPTS=0
RemainAfterExit=no

[Install]
WantedBy=network-online.target
